# includes/supervisor.ini
[program:job_runner_broadcast_queue]
command=/opt/virtualenv/job-runner/bin/manage.py broadcast_queue --settings=job_runner.settings.env.production
user=pyuser
